# photobooth_init.conf
# Photobooth service

description "Photobooth service"
author "pi"

start on runlevel [2345]
stop on runlevel [016]
chdir /home/pi/git/raspberrypi_photobooth
exec python /home/pi/git/raspberrypi_photobooth/PhotoBooth.py
respawn
